# setup/radioman.service
[Unit]
Description=Radioman — Wi-Fi audit daemon
After=network-online.target multi-user.target
Wants=network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/radioman
# 30-second delay gives SSH a window to connect before Python eats all RAM
ExecStartPre=/bin/sleep 30
ExecStart=/opt/radioman/venv/bin/python3 /opt/radioman/radioman.py
Restart=on-failure
RestartSec=60
StandardOutput=journal
StandardError=journal
SyslogIdentifier=radioman

# Give bettercap and airmon-ng the raw socket access they need
AmbientCapabilities=CAP_NET_RAW CAP_NET_ADMIN
CapabilityBoundingSet=CAP_NET_RAW CAP_NET_ADMIN

Environment=PYTHONUNBUFFERED=1
Environment=RADIOMAN_CONFIG=/opt/radioman/radioman.conf

LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
